# waste.rc
#!/bin/ksh

daemon="/usr/local/bin/waste"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
